# etc/linux-upstart/init-headphone.conf
description "Initialize headphone amplifier found in some Clevo laptops"

start on runlevel [2345]

exec /usr/local/sbin/init-headphone
